# Fonera/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

ntpdate es.pool.ntp.org
stty -F /dev/ttyS0 cs8 9600 ignbrk -brkint -icrnl -imaxbel -opost -onlcr -isig -icanon -iexten -echo -echoe -echok -echoctl -echoke noflsh -ixon -crtscts
/srv/www/lee_temp.sh
exit 0
